# agrodiary/deploy/agrodiary.service
[Unit]
Description=AgroDiary - Cuaderno de Campo Digital
Documentation=https://github.com/tu-repo/agrodiary
After=network.target

[Service]
Type=simple
User=agrodiary
Group=agrodiary
WorkingDirectory=/opt/agrodiary
ExecStart=/usr/bin/node /opt/agrodiary/.next/standalone/server.js
Restart=always
RestartSec=5

# Variables de entorno
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=TZ=Europe/Madrid
EnvironmentFile=-/opt/agrodiary/.env.local

# Seguridad
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agrodiary/data /opt/agrodiary/public/uploads /opt/agrodiary/backups

# Límites
LimitNOFILE=65536

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agrodiary

[Install]
WantedBy=multi-user.target
